# notion-daemon.service
[Unit]
Description=Notion Automation Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/Notion_Automator
ExecStart=/path/to/Notion_Automator/venv/bin/python daemon.py
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
